# custom/systemd/system/flatpak-preinstall.service
[Unit]
Description=Preinstall Flatpaks
After=network-online.target
Wants=network-online.target
ConditionPathExists=/usr/bin/dbus-run-session
ConditionPathExists=/usr/bin/flatpak
Documentation=man:flatpak-preinstall(1)

[Service]
Type=oneshot
ExecStart=/usr/bin/dbus-run-session /usr/bin/flatpak preinstall -y --system --noninteractive
RemainAfterExit=true
Restart=on-failure
RestartSec=30

StartLimitIntervalSec=600
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
